# var-backups-nextcloud.mount
[Mount]
What=/dev/vg.nextcloud/nc-backups
Where=/var/backups/nextcloud

[Install]
WantedBy=multi-user.target
