# roles/nextcloud/files/etc/systemd/system/nextcloud-backup.timer
[Unit]
Description=Run Nextcloud backup every hour

[Timer]
Persistent=true
OnCalendar=*-*-* *:00:00

[Install]
WantedBy=timers.target
